# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cpa]
command=/usr/local/bin/CLIProxyAPI -config /etc/cpa/config.yaml
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# 注意：当 GITSTORE_GIT_URL 环境变量设置时，CPA 自动切换到 GitStore 模式，
# 忽略 -config 参数，从 gitstore/config/config.yaml 加载配置。
# 未设置 GITSTORE_GIT_URL 时，使用本地 /etc/cpa/config.yaml。

[program:cpamp]
command=/usr/local/bin/start-cpamp.sh
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sync-data]
command=/usr/local/bin/sync-data.sh
autostart=true
autorestart=true
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
